# supervisord.conf
[supervisord]
nodaemon=true

[program:Xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1280x1024x24
autorestart=true

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -rfbport 5900 -forever -shared -nopw
autorestart=true

[program:app]
command=python main.py
autorestart=true